# system-service/openauto-system.service
# Reference only — install.sh generates the actual service file
# with correct paths for the target system.
#
# Installed to: /etc/systemd/system/openauto-system.service

[Unit]
Description=OpenAuto Prodigy System Manager
Before=openauto-prodigy.service
After=network.target bluetooth.target

[Service]
Type=notify
User=matt
ExecStart=/home/matt/openauto-prodigy/system-service/.venv/bin/python3 /home/matt/openauto-prodigy/system-service/openauto_system.py
ExecStopPost=-/usr/sbin/iptables -t nat -D OUTPUT -p tcp -j OPENAUTO_PROXY
ExecStopPost=-/usr/sbin/iptables -t nat -F OPENAUTO_PROXY
ExecStopPost=-/usr/sbin/iptables -t nat -X OPENAUTO_PROXY
WorkingDirectory=/home/matt/openauto-prodigy/system-service
RuntimeDirectory=openauto
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
