# system-service/openauto-system.service.in
[Unit]
Description=OpenAuto Prodigy System Manager
Before=openauto-prodigy.service
After=network.target bluetooth.target

[Service]
Type=notify
User=root
ExecStart=@@PYTHON_PATH@@ @@SYS_DIR@@/openauto_system.py
ExecStopPost=-/usr/sbin/iptables -t nat -D OUTPUT -p tcp -j OPENAUTO_PROXY
ExecStopPost=-/usr/sbin/iptables -t nat -F OPENAUTO_PROXY
ExecStopPost=-/usr/sbin/iptables -t nat -X OPENAUTO_PROXY
WorkingDirectory=@@SYS_DIR@@
RuntimeDirectory=openauto
Restart=always
RestartSec=2

# Minimal hardening — root daemon, not an enterprise sandbox
PrivateTmp=yes
ProtectHome=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
